# .config/systemd/user/Packages.service
[Unit]
Description=Backup explicitly installed packages to a list

[Service]
Type=oneshot
ExecStart=%h/bin/Packages.sh
